# configs/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:doge]
command=yarn start
directory=/app
autorestart=true
environment=UV_PORT=3000,BARE_PORT=8080